# systemd/gelber_sack.timer
[Unit]
Description=hourly timer for the gelben sack

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
